# docker/scrapy-service/supervisord.conf
[supervisord]
nodaemon=true

[program:scrapy-service]
process_name=%(program_name)s_%(process_num)02d
user=root
command=npm run dev
directory=/app
autostart=true
autorestart=true
numprocs=1
